# etc/seisma/cron.d/fix_aborted
## Fix aborted builds

5 * * * *  seisma  SEISMA_SETTINGS=/etc/seisma/seisma.py PYTHONPATH=/usr/local/src/seisma-server python3.5 -m seisma fix_aborted